# kittysnap.cron
# Every halfhour run job
0,30 * * * * mike_sergio /Users/mike_sergio/code/go/src/github.com/michaelsergio/kittyspy/kittysnap

