# kittysnap.cron
# Every halfhour run job
*/30 5-23 * * * /Users/msergio/code/go/src/github.com/michaelsergio/kittysnap/kittysnap-cron.sh

